# config/systemd/user/ssh-agent-relay.service
[Unit]
Description=1Password SSH Agent Relay
After=network.target

[Service]
ExecStart=%h/.local/bin/ssh-agent-relay.sh
Restart=on-failure
RestartSec=10
TimeoutStartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
